# api/deployment/uat/service.ini
[program:spotify_playlist_analyzer_api_uat]
directory=%(ENV_HOME)s/html/spotify-playlist-analyzer/api/uat
command=%(ENV_HOME)s/html/spotify-playlist-analyzer/api/uat/ENV/bin/uwsgi uwsgi.ini
autostart=yes
autorestart=yes
