# flask-app/deployment/uat/service.ini
[program:spotify_playlist_analyzer_uat]
directory=%(ENV_HOME)s/html/spotify-playlist-analyzer/flask-app/uat
command=%(ENV_HOME)s/html/spotify-playlist-analyzer/flask-app/uat/ENV/bin/uwsgi uwsgi.ini
autostart=yes
autorestart=yes
