# flask-app/deployment/uat/spotify_playlist_analyzer_test.ini
[program:spotify_playlist_analyzer_test]
directory=%(ENV_HOME)s/html/spotify-playlist-analyzer/flask-app/test
command=%(ENV_HOME)s/html/spotify-playlist-analyzer/flask-app/test/ENV/bin/uwsgi uwsgi.ini
autostart=yes
autorestart=yes
